# scraper/cron.d/scraper-cron
0 2 * * * node /usr/src/app/scraper.js >> /usr/src/app/output.log 2>&1
